# debian.default/puavo-load-balancer.upstart
description "puavo-load-balancer"

start on runlevel [23]
stop on shutdown

exec /usr/sbin/puavo-post-ltsp-server-load

respawn
